# templates/default/tyk-gateway.service.erb
[Unit]
Description=Tyk API Gateway

[Service]
LimitNOFILE=<%= @open_file_limit %>
LimitNPROC=<%= @open_file_limit %>
Type=simple
User=root
Group=root
ExecStart=/opt/tyk-gateway/tyk --conf=/opt/tyk-gateway/tyk.conf
Restart=always

[Install]
WantedBy=multi-user.target
